# crontab
# ┌───────────── minute (0 - 59)
# │ ┌───────────── hour (0 - 23)
# │ │ ┌───────────── day of month (1 - 31)
# │ │ │ ┌───────────── month (1 - 12)
# │ │ │ │ ┌───────────── day of week (0 - 6) (Sunday=0)
# │ │ │ │ │
# │ │ │ │ │
* * * * * root python /src/retriever/elos.py && python /src/retriever/fixtures.py && python /src/simulator/sim_season.py >> /var/log/simulator.log 2>&1
